# base/files/etc/systemd/system/networking.service.d/override.conf
# start after container-boot service
[Unit]
After=container-boot.service

